# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:dst-admin]
command=/opt/dst/admin/mod-manager -root /opt/dst/admin -listen 0.0.0.0 -port 8788 -dst-dir /data
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dst-master]
command=/bin/bash -lc 'if [[ "${DST_DISABLE_GAME_PROCESSES:-false}" == "true" ]]; then echo "dst-master disabled until Klei token is configured"; sleep infinity; fi; exec /opt/dst/game/bin64/dontstarve_dedicated_server_nullrenderer_x64 -skip_update_server_mods -ugc_directory /data/ugc_mods -persistent_storage_root /data -conf_dir cluster -cluster Cluster_1 -shard Master'
autorestart=unexpected
startsecs=20
stopwaitsecs=720
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dst-caves]
command=/bin/bash -lc 'if [[ "${DST_DISABLE_GAME_PROCESSES:-false}" == "true" ]]; then echo "dst-caves disabled until Klei token is configured"; sleep infinity; fi; exec /opt/dst/game/bin64/dontstarve_dedicated_server_nullrenderer_x64 -skip_update_server_mods -ugc_directory /data/ugc_mods -persistent_storage_root /data -conf_dir cluster -cluster Cluster_1 -shard Caves'
autorestart=unexpected
startsecs=20
stopwaitsecs=720
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
