# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# 设计：dst-admin 自动启动并守护；dst-master / dst-caves 完全由管理端通过 supervisorctl
# 控制（autostart=false）。这样未配置 token 或未点击启动时，DST 进程就不会运行，
# 也不会出现失败重启循环。

[program:dst-admin]
command=/opt/dst/admin/mod-manager -root /opt/dst/admin -dst-dir /data -listen 0.0.0.0 -port 8788 -supervisor-conf /opt/dst/runtime/supervisord.conf
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dst-master]
command=/opt/dst/game/bin64/dontstarve_dedicated_server_nullrenderer_x64 -skip_update_server_mods -ugc_directory /data/ugc_mods -persistent_storage_root /data -conf_dir cluster -cluster Cluster_1 -shard Master
directory=/opt/dst/game/bin64
autostart=false
autorestart=unexpected
startsecs=20
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dst-caves]
command=/opt/dst/game/bin64/dontstarve_dedicated_server_nullrenderer_x64 -skip_update_server_mods -ugc_directory /data/ugc_mods -persistent_storage_root /data -conf_dir cluster -cluster Cluster_1 -shard Caves
directory=/opt/dst/game/bin64
autostart=false
autorestart=unexpected
startsecs=20
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
